# my_files/luci-app-wifimgr/root/etc/uci-defaults/99-wifimgr-defaults
#!/bin/sh

# Guard: only run once. /etc/config is preserved across sysupgrade,
# so this sentinel survives firmware upgrades — re-running would overwrite
# user-configured passwords, encryption and country settings.
[ "$(uci -q get wifimgr.defaults.applied)" = '1' ] && exit 0

# Set country to CZ only if still at factory default (unset, WORLD, or 00).
# Skip if already configured — preserves country on upgrades from older versions.
_cur=$(uci -q get wireless.radio0.country)
if [ -z "$_cur" ] || [ "$_cur" = "WORLD" ] || [ "$_cur" = "00" ]; then
    uci set wireless.radio0.country='CZ'
    uci set wireless.radio1.country='CZ'
    uci set wireless.radio2.country='CZ'
fi

# On factory interfaces only: set default password + encryption + rename SSID.
# Factory SSIDs are MT76_AP_MLD and OpenWrt (or OpenWrt-Nxxxx).
for iface in $(uci show wireless | grep '=wifi-iface' | cut -d'.' -f2 | cut -d'=' -f1); do
    ssid=$(uci -q get wireless.${iface}.ssid)
    case "$ssid" in
        MT76_AP_MLD|OpenWrt|OpenWrt-*)
            uci set wireless.${iface}.key='12345678'
            radio=$(uci -q get wireless.${iface}.device)
            [ "$radio" = 'radio2' ] && enc='sae' || enc='sae-mixed'
            uci set wireless.${iface}.encryption="$enc"
            [ "$ssid" = 'MT76_AP_MLD' ] && uci set wireless.${iface}.ssid='OpenWrt-MLD'
            ;;
    esac
done

uci commit wireless

# Ensure br-lan has a proper /24 netmask (factory default may leave it as /32)
uci set network.lan.netmask='255.255.255.0'
uci commit network

# Mark as applied — stored in /etc/config/wifimgr (preserved across sysupgrade)
uci set wifimgr.defaults=config
uci set wifimgr.defaults.applied=1
uci commit wifimgr
